# covid.conf
[program:covid_gandaki]
command = /home/covid/deploy/covid_gandaki/gunicorn.sh                    ; Command to start app SUPERVISOR FILE
user = covid                                                          ; User to run as
stdout_logfile = /home/covid/deploy/logs/gunicorn_supervisor.log   ; Where to write log messages
redirect_stderr = true                                                ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8                       ; Set UTF-8 as default encoding